# config/app.socket
[Unit]
Description=wsgi socket

[Socket]
ListenStream=/run/app.sock

[Install]
WantedBy=sockets.target